# arch (old)/install/post-install
## First things to do after Installation

################
## SSH & SUDO ##
################

# Enable SSH & Haveged
	systemctl enable sshd
	systemctl enable haveged

# Set up SSH keys

## General Recommendations
# https://wiki.archlinux.org/index.php/General_recommendations
# Users & Groups
	useradd -m -G wheel -s /bin/bash fileserver
	EDITOR=nano visudo
		fileserver ALL=(ALL) ALL

##############
## Security ##
##############

## Extra security measures to install / configure after initial Arch Installation
# https://wiki.archlinux.org/index.php/Security

## Only allow certain users
	nano /etc/pam.d/su && nano /etc/pam.d/su-l
	# "Uncomment following line to require user to be in wheel group"

## Restrict root login
	passwd --lock root

## Deny root access for ssh
	# https://wiki.archlinux.org/index.php/Secure_Shell#Deny
	nano /etc/ssh/sshd_config
		PermitRootLogin no 

	# Denying root access should not only be done through the SSH daemon, but also using pam
	nano /etc/pam.d/sshd
		auth      required  pam_listfile.so      onerr=succeed item=user sense=deny file=/etc/denysshusers

	nano /etc/denysshusers
		root

	chmod 644 /etc/denysshusers

## Kernel hardening
	# Restricting access to kernel logs
	# Restricting access to kernel pointers in the proc filesystem
	# Disable kexec

## Network & Firewalls
	# Firewall
	pacman -Syu ufw
	ufw default deny
	ufw allow to any port 22 proto tcp # to keep ssh running if you are doing this install over ssh
	ufw enable
	# Set up rules according to ufwrules

	# Kernel parameters
		# TCP / IP stack hardening

	# SSH
		# Set up ssh keys

########################
## Package Management ##
########################

## Synchronize time
	timedatectl set-timezone Europe/Brussels

## File Index & Search
	# Run command to update search index database
	updatedb

## Console Improvements
	# Console prompt - Console Bach prompt
		# Add neofetch prompt
		pacman -Syu neofetch
		nano ~/.bashrc
			if [ -f /usr/bin/neofetch ]; then neofetch; fi

## Autologin
systemctl edit getty@tty1

	[Service]
	ExecStart=
	ExecStart=-/usr/bin/agetty --autologin fileserver --noclear %I $TERM

#################
## Performance ##
#################

	# Change cpupower default settings -> governor to performance
	nano /etc/default/cpupower
		governor='performance'

	systemctl enable cpupower

	# You can measure current cpu frequency with below command
		watch grep \"cpu MHz\" /proc/cpuinfo

	## Turn off CPU mitigations - don't do this on critical systems
	nano /etc/default/grub
		mitigations=off

	## Reduce swappiness
	nano /etc/sysctl.d/99-swappiness.conf
		vm.swappiness=10

	## Increase network performance
	nano /etc/sysctl.d/99-sizereceivequeue.conf
		net.core.netdev_max_backlog = 100000
		net.core.netdev_budget = 50000
		net.core.netdev_budget_usecs = 5000

	nano /etc/sysctl.d/99-increasemaxconnections.conf
		net.core.somaxconn = 1024

	nano /etc/sysctl.d/99-increasememorynetinterface.conf
		net.core.rmem_default = 1048576
		net.core.rmem_max = 16777216
		net.core.wmem_default = 1048576
		net.core.wmem_max = 16777216
		net.core.optmem_max = 65536
		net.ipv4.tcp_rmem = 4096 1048576 2097152
		net.ipv4.tcp_wmem = 4096 65536 16777216
		net.ipv4.udp_rmem_min = 8192
		net.ipv4.udp_wmem_min = 8192

	nano /etc/sysctl.d/99-tcpfastopen.conf
		net.ipv4.tcp_fastopen = 3

	nano /etc/sysctl.d/99-bbr.conf
		net.core.default_qdisc = fq
		net.ipv4.tcp_congestion_control = bbr
