# SystemdFiles/autoupdate.service
[Unit]
Description = ArchServer Auto Update Service
OnFailure=unit-failure-notification@%n

[Service]
User=fileserver
Group=fileserver
Type=oneshot
ExecStart=/home/fileserver/Applications/autoupdate/autoupdate.sh
