# arch (old)/systemdfiles/clonezillabackup.service
[Unit]
Description = CloneZilla Backup Script

[Service]
Type=oneshot
ExecStart=/home/fileserver/ArchServer/scripts/clonezillabackup.sh