# arch (old)/systemdfiles/deletelatenight.service
[Unit]
Description = Delete Late Night Shows Service

[Service]
Type=oneshot
ExecStart=/home/fileserver/ArchServer/scripts/deletelatenight.sh
